# services.d/nodebook/run
#!/usr/bin/with-contenv sh

echo "------------------------------------------------------------"
echo "------------------------------------------------------------"
echo "nodebook init at - $(date)"
echo "------------------------------------------------------------"
echo "------------------------------------------------------------"

/app/nodebook-linux web --bindaddress=0.0.0.0 --port="$NODEBOOK_PORT" --docker /notebooks
